# host/mining-appliance/bdag-node-child-guard.timer
[Unit]
Description=Run BlockDAG node child liveness guard

[Timer]
OnBootSec=2min
OnUnitActiveSec=60s
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
