# ops/systemd/user-bdag-chain-presync.service
[Unit]
Description=BlockDAG live chain pre-sync
After=default.target

[Service]
Type=oneshot
WorkingDirectory=/home/jeremy/blockdag-mining-pool/stack
Environment=BDAG_PROJECT_ROOT=/home/jeremy/blockdag-mining-pool/stack
Environment=BDAG_PRESYNC_ONE_NODE=1
EnvironmentFile=-/home/jeremy/blockdag-mining-pool/stack/ops/runtime/ops.env
Nice=19
IOSchedulingClass=idle
CPUWeight=10
IOWeight=10
ExecStart=/home/jeremy/blockdag-mining-pool/stack/ops/chain-presync.sh
