# ops/systemd/user-bdag-fastsnap-seed.service
[Unit]
Description=BlockDAG public FastSnap seed export
After=default.target

[Service]
Type=oneshot
Environment=BDAG_FASTSNAP_REQUIRE_BOTH_BACKENDS_FOR_VERIFY=0
Environment=BDAG_FASTSNAP_REQUIRE_EXPORT_BACKEND_FRESH=1
Environment=BDAG_FASTSNAP_MAX_EXPORT_BACKEND_LAG=1000
EnvironmentFile=-%h/.config/bdag-fastsnap-seed.env
Nice=19
IOSchedulingClass=idle
CPUWeight=10
IOWeight=10
ExecStart=/usr/bin/env bash -lc 'cd "${BDAG_PROJECT_ROOT:?set BDAG_PROJECT_ROOT in ~/.config/bdag-fastsnap-seed.env}" && exec ./ops/build-fastsnap-seed.sh'
