# systemd/miner-blockchain-backup-retry.service
[Unit]
Description=Retry pending BDAG Miner blockchain backup once node is healthy
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/miner/pool-stack-docker
Environment=KEEP=3
Environment=BACKUP_ROOT=/opt/miner/backups/blockchain
Environment=STATE_ROOT=/opt/miner/backup-state
Environment=LOG_DIR=/opt/miner/logs
Environment=POSTGRES_USER=bdag_pool
Environment=POSTGRES_DB=bdagpool
ExecStart=/opt/miner/pool-stack-docker/scripts/miner-blockchain-cold-backup.sh --retry-pending
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
